# party-game.service
[Unit]
Description=Party Game Server
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/party-game
ExecStart=/home/pi/party-game/start.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target 